# deploy/systemd/cyberdeck-backend.service
[Unit]
Description=CyberDeck Backend
After=network.target redis-cyberdeck.service
Requires=redis-cyberdeck.service

[Service]
Type=simple
WorkingDirectory=/opt/cyberdeck/backend
Environment=REDIS_URL=redis://127.0.0.1:6379
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port 8000
Restart=on-failure
User=cyberdeck

[Install]
WantedBy=multi-user.target
