# scripts/mesh.service
[Unit]
Description=CyberDeck AirGapped Mesh Node Daemon
After=network.target
Wants=network-online.target

[Service]
Type=simple
# Path to the python script
ExecStart=/usr/bin/env python3 /opt/cyberdeck/scripts/mesh_node.py
Restart=always
RestartSec=5
# Run as root to allow 'ip' commands
User=root

# Environment variables can be overridden here
Environment="MESH_SECRET=GHOSTWIRE_ALPHA_7"
Environment="MESH_IFACE=eth0"
Environment="MESH_BEACON_PORT=5005"
Environment="MESH_CONTROL_PORT=5006"
Environment="MESH_API_PORT=5007"

[Install]
WantedBy=multi-user.target
